# infra/systemd/peilv-dispatch.timer
[Unit]
Description=Dispatch Peilv automation tasks every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1s
Unit=peilv-dispatch.service

[Install]
WantedBy=timers.target
